# overlays/System/Library/LaunchDaemons/com.apple.SystemConfiguration.KernelEventMonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.apple.SystemConfiguration.KernelEventMonitor.plist

  Boot-time launchd plist for KernelEventMonitor (freebsd-launchd-mach
  SystemConfiguration kernel-event monitor). Apple runs this as a
  configd plugin inside the configd process; we run it standalone
  because this repo's configd has no plugin loader and stock FreeBSD
  has no EVFILT_MACHPORT (so configd's single mach_msg loop can't
  multiplex a Mach port and the route socket). Folding it into configd
  is tracked as a kernel-track item.

  It reads the PF_ROUTE socket for RTM_IFINFO link-state changes and
  publishes State:/Network/Interface/<if>/Link = {Active:…} into
  configd's SCDynamicStore. ipconfigd watches that key and starts DHCP
  when a NIC's link comes up — replacing the removed hwregd attach
  trigger. RunAtLoad + KeepAlive so it starts at PID-1 launchd boot and
  respawns on crash. StandardOutPath / StandardErrorPath route the log
  to /var/log/KernelEventMonitor.stderr — same convention configd /
  ipconfigd use; the CI boot test dumps the file post-boot.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.apple.SystemConfiguration.KernelEventMonitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/sbin/KernelEventMonitor</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/KernelEventMonitor.stderr</string>
    <key>StandardErrorPath</key>
    <string>/var/log/KernelEventMonitor.stderr</string>
</dict>
</plist>
