# overlays/System/Library/LaunchDaemons/com.apple.hostnamed.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.hostnamed</string>

	<!-- Persistent libdispatch event loop (iter 3c). Runs initial
	     synthesize + sethostname at boot, then enters dispatch_main()
	     waiting on SCDynamicStore notifications (DHCP Option_12 +
	     IPv4 changes), SCPreferencesSetCallback (ComputerName edits),
	     and a SIGHUP dispatch source (operator/CI re-synth nudge,
	     also the mDNS PTR appearance trigger since that's not an
	     SCDS event). Refreshes the 5-tier precedence chain (kenv >
	     SCPrefs > DHCP > mDNS > synthesis) on every callback; the
	     refresh memoizes the last-published value so no-op rebuilds
	     are skipped.

	     KeepAlive=true: if the daemon ever crashes, launchd restarts
	     it — the hostname is system-of-record state, can't be
	     ephemeral. SIGTERM (launchctl unload) is the clean-shutdown
	     path; the dispatch SIGTERM source exits cleanly.

	     PAM port iter 4 (issue #99): RunAtLoad RESTORED. The
	     original iter-1 drop was specifically because the additional
	     RunAtLoad plist accelerated launchd's boot-time plist
	     dispatch enough to let getty/login race ahead of pam_xdg's
	     /var/run/xdg setup, breaking pam_open_session and wedging
	     root login. PAM port iter 3 (PR #98) dropped pam_xdg from
	     the overlay /etc/pam.d/* entirely; nothing references it
	     anymore. The race is structurally resolved.

	     With RunAtLoad=true, hostnamed runs at boot before getty
	     prints its banner; the initial refresh_hostname() calls
	     sethostname() FIRST (before SCDS publishes do their CF
	     allocation), so the kernel hostname is set before getty
	     can print its banner. The synthesized hostname appears in
	     "FreeBSD/amd64 (HOSTNAME) (console)" instead of the
	     unset-placeholder "Amnesiac".

	     NOTE: com.apple.syslogd.plist has a separate, unrelated
	     reason its RunAtLoad is dropped (the launch_msg LAUNCH_KEY_CHECKIN
	     Mach IPC hang from the mach.ko interruptibility gap). That
	     issue isn't fixed by the PAM port — syslogd stays dropped.

	     Plan: https://pkgdemon.github.io/freebsd-hostnamed-plan.html
	     Issue: #63 (iter 1), #99 (RunAtLoad), iter 3c (persistent loop) -->
	<key>ProgramArguments</key>
	<array>
		<string>/usr/sbin/hostnamed</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>

	<!-- Capture stderr for post-mortem + the HOSTNAMED-OK marker
	     visibility from the on-image test runner. -->
	<key>StandardErrorPath</key>
	<string>/var/log/hostnamed.stderr</string>
	<key>StandardOutPath</key>
	<string>/var/log/hostnamed.stderr</string>

	<!-- System session — runs as root before any user logs in. Same
	     scope as configd / hwregd / notifyd / mDNSResponder. -->
	<key>LimitLoadToSessionType</key>
	<string>System</string>
</dict>
</plist>
