# overlays/System/Library/LaunchDaemons/org.freebsd.hwmatch-test.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  org.freebsd.hwmatch-test.plist

  Phase K fixture for the launchd HardwareMatch plist key. launchd
  PID 1 imports this plist at boot; parsing its HardwareMatch <array>
  and HardwareMatchAction <string> exercises hwmatch_new_from_obj()
  and lands the job on launchd's global HardwareMatch list.

  iter 2 wired launchd to hwregd's device-event notify bus: at boot,
  when hwregd checks its Mach service in, launchd subscribes; on a
  device arrival whose BSD name equals a job's HardwareMatch `name`
  criterion, launchd dispatches the job. iter 2 matches the device
  name only (the SUBSCRIBE event line carries no class / driver).

  The criterion below names a device that no QEMU CI guest has, so
  the job does not fire during the boot test — the CI gate for iter 2
  is the launchd<->hwregd subscription round-trip (the "subscription
  to hwregd device events confirmed" console line), not this job
  firing. /usr/bin/true is a harmless placeholder program; the job
  has no RunAtLoad / KeepAlive so it only ever runs on a matching
  device arrival.

  Plan: pkgdemon.github.io/freebsd-hardware-registry-iokit-plan.html §4
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.freebsd.hwmatch-test</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/true</string>
    </array>
    <key>HardwareMatch</key>
    <array>
        <dict>
            <key>name</key>
            <string>hwmatchtest0</string>
        </dict>
    </array>
    <key>HardwareMatchAction</key>
    <string>StartOnArrival</string>
</dict>
</plist>
