# overlays/System/Library/LaunchDaemons/org.freebsd.hwmatch-test.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  org.freebsd.hwmatch-test.plist

  Phase K iteration 1 fixture for the launchd HardwareMatch plist
  key. launchd PID 1 imports this plist at boot; parsing its
  HardwareMatch <array> and HardwareMatchAction <string> exercises
  the new hwmatch_new_from_obj() path in src/launchd/src/core.c.

  iter 1 only parses and stores the criteria — there is no
  device-arrival trigger yet, and the job has no RunAtLoad /
  KeepAlive, so /usr/bin/true is never actually run. The boot test
  staying green is the iter-1 verification: a HardwareMatch parsing
  bug would crash launchd PID 1 and fail every CI marker.

  The hwreg_watch RPC client that registers these criteria with
  hwregd and fires the job on a matching device arrival lands in
  iter 2. The criteria keys mirror hwregd's watch filter — name /
  class / driver (hwreg.defs hwreg_watch).

  Plan: pkgdemon.github.io/freebsd-hardware-registry-iokit-plan.html §4
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.freebsd.hwmatch-test</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/true</string>
    </array>
    <key>HardwareMatch</key>
    <array>
        <dict>
            <key>class</key>
            <string>cpu</string>
        </dict>
    </array>
    <key>HardwareMatchAction</key>
    <string>StartOnArrival</string>
</dict>
</plist>
