# overlays/System/Library/LaunchDaemons/org.freebsd.coredump-off.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  org.freebsd.coredump-off.plist

  Disables process core dumps system-wide at boot (sysctl
  kern.coredump=0). Workaround: a crashing daemon whose fdtable holds
  mach.ko entries trips a FreeBSD kernel page fault in
  kern_proc_filedesc_out during coredump. Until mach.ko's fd ops are
  coredump-safe, dumps stay off.

  This was previously done by ramdisk/init.sh before the unionfs
  pivot. With the disk-image boot there is no init.sh and launchd
  bypasses rc.d (so /etc/sysctl.conf is never applied) — so a one-shot
  RunAtLoad launchd job carries it. kern.coredump is not a loader
  tunable, so this is the earliest available hook.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.freebsd.coredump-off</string>
    <key>ProgramArguments</key>
    <array>
        <string>/sbin/sysctl</string>
        <string>kern.coredump=0</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
