# overlays/System/Library/LaunchDaemons/org.freebsd.hwregd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  org.freebsd.hwregd.plist

  Boot-time launchd plist for the freebsd-launchd-mach hardware
  registry daemon. Replaces FreeBSD-devd's role for kldload-on-
  nomatch + hardware event fanout (devd was removed at commit
  88694f0; hwregd is the structured replacement, plan at
  pkgdemon.github.io/freebsd-hardware-registry-iokit-plan.html).

  iter 2: daemon reads /dev/devctl and, on a nomatch event,
  reports which module linker.hints says claims the device.
  RunAtLoad=true so it starts at PID-1 launchd boot; KeepAlive=true
  so any crash gets respawned. StandardErrorPath/StandardOutPath
  route its log to /var/log/hwregd.stderr — the path the CI boot
  test (run.sh) dumps post-boot for smoke verification.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.freebsd.hwregd</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/sbin/hwregd</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>MachServices</key>
    <dict>
        <key>org.freebsd.hwregd</key>
        <true/>
    </dict>
    <key>StandardOutPath</key>
    <string>/var/log/hwregd.stderr</string>
    <key>StandardErrorPath</key>
    <string>/var/log/hwregd.stderr</string>
</dict>
</plist>
